# deploy/systemd/chloe-hindsight-nightly.service
[Unit]
Description=Chloe Full Hindsight Cycle (Reflection + Tuner + Dream)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/Chloe-alpha
Environment="PYTHONPATH=/root/Chloe-alpha"
ExecStart=/root/Chloe-alpha/venv/bin/python -m tools.hindsight_cycle full
