# deploy/systemd/chloe-hindsight-nightly.timer
[Unit]
Description=Run Chloe Full Hindsight nightly

[Timer]
OnCalendar=03:30
Persistent=true

[Install]
WantedBy=timers.target
