# systemd/chloe-exploration-accelerator.service.example
[Unit]
Description=Chloe Exploration Accelerator (advisory overrides)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/Chloe-alpha
EnvironmentFile=-/root/Chloe-alpha/systemd/chloe.env
ExecStart=/bin/bash -lc 'source venv/bin/activate && python3 -m tools.run_exploration_accelerator'

[Install]
WantedBy=multi-user.target

