# systemd/chloe-shadow-promotion-rollback.timer.example
[Unit]
Description=Run Chloe Shadow Promotion Rollback every 10 minutes (offset)

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target

