# systemd/chloe-reflect-performance.timer.example
[Unit]
Description=Run Chloe GPT Reflection (Performance) periodically

[Timer]
# Default: every 10 minutes. Adjust to your preference.
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target

