# systemd/chloe-reflect-performance.service.example
[Unit]
Description=Chloe GPT Reflection (Performance) - oneshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/Chloe-alpha
Environment=PYTHONUNBUFFERED=1
# NOTE: systemd env files do NOT support `export ...`.
# Create /root/Chloe-alpha/systemd/chloe.env (see systemd/chloe.env.example).
EnvironmentFile=-/root/Chloe-alpha/systemd/chloe.env
ExecStart=/bin/bash -lc 'source venv/bin/activate && python3 -m engine_alpha.reflect.gpt_reflection_runner'

[Install]
WantedBy=multi-user.target

