# systemd/chloe-reflection-v4.timer.example
[Unit]
Description=Run Chloe GPT Reflection v4 (Narrative) nightly

[Timer]
OnCalendar=*-*-* 03:15:00
AccuracySec=10min
Persistent=true

[Install]
WantedBy=timers.target

