# systemd/nightly_research.timer.example
[Unit]
Description=Run Chloe Nightly Research daily

[Timer]
OnCalendar=*-*-* 03:05:00
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target

